# deploy/systemd/overnet-irc.service
[Unit]
Description=Overnet IRC Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/overnet/overnet-irc.env
ExecStart=/usr/bin/env perl /opt/overnet/overnet-program-irc/bin/overnet-irc-service.pl \
  --adapter-id ${OVERNET_IRC_ADAPTER_ID} \
  --network ${OVERNET_IRC_NETWORK} \
  --listen-host ${OVERNET_IRC_LISTEN_HOST} \
  --listen-port ${OVERNET_IRC_LISTEN_PORT} \
  --server-name ${OVERNET_IRC_SERVER_NAME} \
  --signing-key-file ${OVERNET_IRC_SIGNING_KEY_FILE} \
  --group-host ${OVERNET_IRC_GROUP_HOST} \
  --authority-relay-url ${OVERNET_IRC_AUTHORITY_RELAY_URL} \
  --health-file ${OVERNET_IRC_HEALTH_FILE} \
  --log-file ${OVERNET_IRC_LOG_FILE} \
  --tls \
  --tls-cert-chain-file ${OVERNET_IRC_TLS_CERT_CHAIN_FILE} \
  --tls-private-key-file ${OVERNET_IRC_TLS_PRIVATE_KEY_FILE}
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
